# patches/etc/cron.daily/backup
#!/bin/bash

touch /home/ivan/backup-job-completed
source $NC_ROOT/Backup/cron/daily
