# patches/etc/cron.daily/backup
#!/bin/bash

sudo -u ivan /bin/bash -c 'source $HOME/.profile && cd $NC_BACKUP/cron && ./daily'
