# patches/etc/cron.daily/backup
#!/bin/bash

sudo -u ivan /bin/bash -c 'source $HOME/.profile; $NC_BACKUP/cron/daily'
